# debian/ciborium.conf
description "ciborium manages storage devices"
author "Sergio Schvezov <sergio.schvezov@canonical.com>"

respawn 

start on started dbus
stop on stopping dbus

exec /usr/lib/ciborium/ciborium
